# 環境設定やリリース手順(manusで作成)/gmail-relay-controller.service
[Unit]
Description=Gmail Alert Relay Controller Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
ExecStart=/usr/bin/python3 /home/pi/gmail_relay_controller.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# リソース制限
MemoryLimit=256M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
